# etc/profile.d/kiosk.sh
#!/bin/sh

sudo xinit /usr/bin/startkiosk.sh
